# DIYbyt-Client/systemd/diybyt-display.service
[Unit]
Description=DIYbyt Display Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/DIYbyt

# Environment variables
Environment=PYTHONUNBUFFERED=1
Environment=DIYBYT_SERVER_URL=${SERVER_URL}
Environment=DIYBYT_PROGRAMS_PATH=/opt/DIYbyt/star_programs
Environment=PATH=/opt/DIYbyt/venv/bin:$PATH
# Use exact path where module is found
Environment=PYTHONPATH=/home/ycsgc/DIYbyt/rgb-led-matrix-library/rpi-rgb-led-matrix/bindings/python

# Matrix Configuration
Environment=DIYBYT_MATRIX_ROWS=32
Environment=DIYBYT_MATRIX_COLS=64
Environment=DIYBYT_GPIO_SLOWDOWN=2
Environment=DIYBYT_DISABLE_HARDWARE_PULSING=true

# Service execution
ExecStart=/opt/DIYbyt/venv/bin/python /usr/local/bin/diybyt-display

# Restart configuration
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=append:/var/log/diybyt/display.log
StandardError=append:/var/log/diybyt/display.log

[Install]
WantedBy=multi-user.target