# DIYbyt-Client/systemd/diybyt-display.service
[Unit]
Description=DIYbyt Display Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Restart=always
RestartSec=5

# Server Configuration
Environment=DIYBYT_SERVER_URL=http://localhost

# LED Matrix Configuration
Environment=DIYBYT_MATRIX_ROWS=32
Environment=DIYBYT_MATRIX_COLS=64
Environment=DIYBYT_GPIO_SLOWDOWN=4
Environment=DIYBYT_DISABLE_HARDWARE_PULSING=true

# Optional LED Matrix Parameters (uncomment and modify as needed)
#Environment=DIYBYT_MATRIX_CHAIN_LENGTH=1
#Environment=DIYBYT_MATRIX_PARALLEL=1
#Environment=DIYBYT_MATRIX_BRIGHTNESS=100

# Python Path Configuration
Environment=PYTHONPATH=/opt/DIYbyt:/opt/rpi-rgb-led-matrix/bindings/python

# Execute the script
ExecStart=/usr/bin/python3 /opt/DIYbyt/display.py

[Install]
WantedBy=multi-user.target