# DIYbyt-Server/systemd/diybyt-renderer.service
[Unit]
Description=DIYbyt Renderer Service
After=network.target
Wants=diybyt-sync.service

[Service]
Type=simple
User=ycsgc
Group=ycsgc
WorkingDirectory=/opt/DIYbyt/render
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/DIYbyt/render/venv/bin/python src/components/ProgramManager/pixlet_renderer.py

# Restart configuration
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=append:/var/log/diybyt/renderer.log
StandardError=append:/var/log/diybyt/renderer.log

[Install]
WantedBy=multi-user.target