# DIYbyt-Server/systemd/diybyt-renderer.service
[Unit]
Description=DIYbyt Renderer Service
After=network.target
Wants=diybyt-sync.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/DIYbyt
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/DIYbyt/components/ProgramManager/pixlet_renderer.py

# Restart configuration
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=append:/var/log/diybyt/renderer.log
StandardError=append:/var/log/diybyt/renderer.log

[Install]
WantedBy=multi-user.target