# DIYbyt-Server/systemd/diybyt-renderer.service
[Unit]
Description=DIYbyt Renderer Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=%i
Environment=DIYBYT_HOME=/opt/diybyt
Environment=VIRTUAL_ENV=${DIYBYT_HOME}/venv
Environment=PATH=${VIRTUAL_ENV}/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=${DIYBYT_HOME}
WorkingDirectory=${DIYBYT_HOME}/render
ExecStart=${VIRTUAL_ENV}/bin/python3 ${DIYBYT_HOME}/render/pixlet_renderer.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target