# DIYbyt-Server/systemd/diybyt-renderer.service
[Unit]
Description=DIYbyt Renderer Service
After=network.target

[Service]
Type=simple
User=ycsgc
Group=ycsgc
WorkingDirectory=/opt/DIYbyt/render/src/components/ProgramManager
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/DIYbyt/render/src

ExecStart=/opt/DIYbyt/render/venv/bin/python pixlet_renderer.py

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=append:/var/log/diybyt/renderer.log
StandardError=append:/var/log/diybyt/renderer.log

# Give the service time to gracefully shut down
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target