# DIYbyt-Sync/diybyt-sync.service
[Unit]
Description=DIYbyt Star Programs Sync Service
After=network.target

[Service]
Type=simple
User=ycsgc
Group=ycsgc
Environment=DIYBYT_SERVER_URL=http://100.96.55.97:3001
Environment=DIYBYT_PROGRAMS_PATH=/opt/DIYbyt/star_programs
Environment=DIYBYT_SYNC_INTERVAL=5
ExecStart=/usr/local/bin/diybyt-sync
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target