# DIYbyt-Sync/diybyt-sync.service
[Unit]
Description=DIYbyt Sync Service
After=network.target

[Service]
Type=simple
User=ycsgc
Group=ycsgc
WorkingDirectory=/opt/DIYbyt

# Environment variables
Environment=PYTHONUNBUFFERED=1
Environment=DIYBYT_SERVER_URL=http://localhost:3001
Environment=DIYBYT_PROGRAMS_PATH=/opt/DIYbyt/star_programs
Environment=DIYBYT_SYNC_INTERVAL=5

# Service execution
ExecStart=/usr/bin/python3 /usr/local/bin/diybyt-sync

# Restart configuration
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=append:/var/log/diybyt/sync.log
StandardError=append:/var/log/diybyt/sync.log

[Install]
WantedBy=multi-user.target